# api_manager.service
# File path: /etc/systemd/system/api_manager.service
# sudo systemctl daemon-reload
# sudo systemctl start api_manager
# sudo systemctl enable api_manager
# sudo systemctl status api_manager
[Unit]
Description=API Manager Service
After=network.target

[Service]
User=ubuntu  # 注意修改用户名
WorkingDirectory=/path/to/api_manager  # 注意修改路径
ExecStart=/path/to/api_manager/api_manager --port 3000 --log-dir /path/to/api_manager/logs  # 注意修改路径和端口号
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

